# packaging/debian/supervisor/gosa-dbus.conf
[program:gosa-dbus]
command=/usr/sbin/gosa-dbus
stopsignal=INT
