# packaging/debian/supervisor/gosa-client.conf
[program:gosa-client]
command=/usr/sbin/gosa-client
user=gosa
stopsignal=INT
