# packaging/debian/supervisor/gosa.conf
[program:gosa]
command=/usr/sbin/gosa
user=gosa
stopsignal=INT
