# packaging/debian/gosa-client.service
[Unit]
Description=GOsa client

[Service]
ExecStart=/usr/sbin/gosa-client
Restart=on-failure
User=gosa

[Install]
WantedBy=multi-user.target
