# packaging/debian/gosa-dbus.postinst
#!/bin/sh

set -e

case "$1" in
  configure)

    # Ask the bus to reload the config file
    invoke-rc.d dbus force-reload || true
    mkdir /etc/gosa/shell.d 2>/dev/null || /bin/true
  ;;
  abort-upgrade|abort-remove|abort-deconfigure)
  ;;
  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
  ;;
esac

#DEBHELPER#
