# packaging/debian/gosa-dbus.service
[Unit]
Description=GOsa DBUS agent

[Service]
ExecStart=/usr/sbin/gosa-dbus
Restart=on-failure

[Install]
WantedBy=multi-user.target
