# packaging/debian/gosa.service
[Unit]
Description=GOsa client

[Service]
ExecStart=/usr/sbin/gosa
Restart=on-failure
User=gosa

[Install]
WantedBy=multi-user.target
